# srcs/nginx/supervisord.conf
[supervisord]
nodaemon=true                 

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"

[include]
files = /etc/supervisor/conf.d/*.conf
